# deploy/systemd/aragora.service.d/secrets.conf
[Service]
# Enable AWS Secrets Manager for loading sensitive configuration
Environment="ARAGORA_USE_SECRETS_MANAGER=true"
