# deploy/aragora-watchdog.service
[Unit]
Description=Aragora Health Watchdog
After=aragora.service aragora-staging.service cloudflared.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/aragora/scripts/healthcheck-watchdog.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
